# templates/limesurvey.timer.j2
[Unit]
Description=Run cronjob for limesurvey at {{ limesurvey_path }}

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
